# 404.socket
[Socket]
ListenStream = 88
BindIPv6Only = both

[Install]
WantedBy = sockets.target
