# 404.service
[Unit]
Description = 404 micro-service

[Service]
ExecStart = /home/dgodd/workspace/systemdsocket/systemdsocket
User = dgodd
Group = dgodd
